# src/main/resources/archetype-resources/docker/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:nginx]
command=nginx

[program:sshd]
command=/usr/sbin/sshd -D

[program:${applicationName}]
command=java -jar ${applicationName}.jar
user=${applicationName}
directory=/home/${applicationName}/
